# etc/cron.hourly/cron-checkdns
# run this crontab entry every 5 minutes
*/5 * * * * /openvpn/checkdns.sh